# usr/lib/systemd/system/ansible-pull.target
[Unit]
Description=GitOps with ansible-pull

[Install]
WantedBy=default.target
